# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

exec /etc/openvpn/start.sh
